# sysadmin/vibevoice/com.bumpyclock.vibevoice.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bumpyclock.vibevoice</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/.venv/bin/python</string>
    <string>__PROJECT_DIR__/run_server.py</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>7781</string>
    <string>--startup-timeout</string>
    <string>1800</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>HF_HOME</key>
    <string>__HOME__/.local/share/vibevoice/huggingface</string>
    <key>HF_HUB_OFFLINE</key>
    <string>1</string>
    <key>TOKENIZERS_PARALLELISM</key>
    <string>false</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/VibeVoice/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/VibeVoice/stderr.log</string>
</dict>
</plist>
